# 05-infra/systemd/panai-model-sync.service
[Unit]
Description=ArchIToken PanAI local model sync
After=network-online.target

[Service]
Type=simple
WorkingDirectory=/home/insome/dev/insomeos
ExecStart=/home/insome/.bun/bin/bun /home/insome/dev/insomeos/tools/panai_model_sync.mjs --watch --interval 60
Restart=always
RestartSec=5

[Install]
WantedBy=default.target
